# systemd/userbackups.timer
[Unit]
Description = Run hourly backups of user directories

[Timer]
OnCalendar = 00/8:00
Persistent = true

[Install]
WantedBy = timers.target
