# systemd/userbackups.service
[Unit]
Description = Execute user directory backups

[Service]
Type = oneshot
ExecStart = /usr/bin/systemd-inhibit --who=root --mode=block /home/jason/bin/bash-backups/scripts/userbackups
User = root
Group = root

[Install]
WantedBy = multi-user.target
